# supervisord.conf

[supervisord]
nodaemon=true

[program:entrypoint]
command=/bin/bash /entrypoint.sh couchdb

[program:setup]
command=/bin/bash /entrypoint.sh setup
